# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

# Control socket so `docker exec <c> supervisorctl status` works for debugging.
[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

# Each third-party service is launched through `/usr/bin/env -u …` so the
# backend-only secrets are REMOVED from its environment (the single shared
# container env from env_file would otherwise leak them — the exact
# JWT_SECRET/LIVEKIT_*/ADMIN_PASSWORD set that main's compose scoped out of OME).
# OME still inherits OME_API_TOKEN/OME_WEBHOOK_SECRET, which it legitimately
# reads. NOTE: add any NEW backend-only secret to the -u list on all four.
#
# Start order (priority, low first): valkey → backend → ome/livekit → caddy, so
# the admission webhook (backend) is up before OME starts accepting ingests.

[program:valkey]
command=/usr/bin/env -u JWT_SECRET -u ADMIN_PASSWORD -u LIVEKIT_API_KEY -u LIVEKIT_API_SECRET /usr/local/bin/valkey-server --save "" --loglevel warning
user=valkey
autostart=true
autorestart=true
priority=10
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:backend]
command=/usr/local/bin/stream-backend
directory=/app
user=app
autostart=true
autorestart=true
priority=20
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# -c origin_conf matches the base image's default Cmd; the image ships
# origin_conf/Server.xml (no conf/ dir), so omitting it crash-loops OME.
[program:ome]
command=/usr/bin/env -u JWT_SECRET -u ADMIN_PASSWORD -u LIVEKIT_API_KEY -u LIVEKIT_API_SECRET /opt/ovenmediaengine/bin/ome_launcher.sh -c origin_conf
directory=/opt/ovenmediaengine/bin
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:livekit]
command=/usr/bin/env -u JWT_SECRET -u ADMIN_PASSWORD -u LIVEKIT_API_KEY -u LIVEKIT_API_SECRET /usr/local/bin/livekit-server --config /livekit.yaml
autostart=true
autorestart=true
priority=30
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
command=/usr/bin/env -u JWT_SECRET -u ADMIN_PASSWORD -u LIVEKIT_API_KEY -u LIVEKIT_API_SECRET /usr/local/bin/caddy run --config /etc/caddy/Caddyfile
autostart=true
autorestart=true
priority=40
startsecs=2
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
